# deployment/systemd/pms-backend.service
[Unit]
Description=YEA PMS Backend (Gunicorn)
After=network.target postgresql.service redis.service

[Service]
Type=notify
User=deploy
Group=deploy
WorkingDirectory=/var/www/YEA/PMS/pms-backend
Environment="PATH=/var/www/YEA/PMS/pms-backend/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=core.settings"

# Load environment variables from production env file
EnvironmentFile=/var/www/YEA/PMS/pms-backend/.env.production

# Gunicorn configuration
ExecStart=/var/www/YEA/PMS/pms-backend/venv/bin/gunicorn \
    --workers 4 \
    --worker-class sync \
    --worker-tmp-dir /dev/shm \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 120 \
    --bind unix:/var/www/YEA/PMS/pms-backend/gunicorn.sock \
    --access-logfile /var/log/pms-backend/access.log \
    --error-logfile /var/log/pms-backend/error.log \
    --log-level info \
    core.wsgi:application

# Restart policy
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pms-backend

[Install]
WantedBy=multi-user.target
